# debian/freifunk/freifunk-www-counter/freifunk-www-counter/freifunk-www-counter.cron.d
*/5 * * * * www-data /usr/share/freifunk-www-counter/counter_update.py /var/www/map/nodes.json > /var/www/counter/counter_image.svg
