# debian/freifunk/freifunk-alfred/freifunk-alfred/debian/freifunk-alfred.cron.d
*/5 * * * * root /usr/share/freifunk-alfred/print_map.sh -c - | alfred -s 64
*/5 * * * * root /usr/share/freifunk-alfred/print_service.sh -c - | alfred -s 91

