# debian/freifunk/freifunk-www-status/freifunk-www-status-1/debian/freifunk-www-status.cron.d
*/5 * * * * www-data /usr/share/freifunk-www-status/status_page_create.sh > /var/www/index.html
